# launchd/com.okdoggy.edit-master-wiki.raw-hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.okdoggy.edit-master-wiki.raw-hourly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/doggy/2026/Edit_Master_wiki/scripts/collect_hourly_raw.py</string>
    <string>--max-items</string>
    <string>3</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/doggy/2026/Edit_Master_wiki</string>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/doggy/2026/Edit_Master_wiki/.omx/logs/hourly_raw_collector.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/doggy/2026/Edit_Master_wiki/.omx/logs/hourly_raw_collector.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONIOENCODING</key>
    <string>utf-8</string>
  </dict>
</dict>
</plist>
